# .config/systemd/user/docker-cleanup.service
[Unit]
Description=Cleanup Docker 

[Service]
ExecStart=/usr/bin/docker system prune --force --all --volumes

[Install]
WantedBy=default.target
