# .config/systemd/user/docker-cleanup.timer
[Unit]
Description=Run docker-cleanup

[Timer]
Persistent=true
OnCalendar=monthly

[Install]
WantedBy=timers.target
