# .config/systemd/user/fusion.service
[Unit]
Description=Fusion dashboard
Documentation=https://github.com/Runfusion/Fusion
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=-/usr/bin/docker rm -f fusion
ExecStart=/usr/bin/docker run --rm \
  --name fusion \
  --network host \
  --volume %h/Sources/Fusion:/project \
  --env ANTHROPIC_API_KEY \
  --env OPENAI_API_KEY \
  --env OPENROUTER_API_KEY \
  --env GITHUB_TOKEN \
  --env FUSION_DASHBOARD_TOKEN \
  runfusion/fusion:latest dashboard --host 0.0.0.0 --port 4040
ExecStop=/usr/bin/docker stop fusion
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
