# .config/systemd/user/fusion.service
[Unit]
Description=Fusion dashboard
Documentation=https://github.com/Runfusion/Fusion
After=network-online.target
Wants=network-online.target

[Service]
Environment=NODE_LLAMA_CPP_SKIP_DOWNLOAD_GPU=1
Environment=GIT_CEILING_DIRECTORIES=%h/Sources/Fusion
Environment="GIT_SSH_COMMAND=ssh -F %h/.fusion/ssh_config"
Environment=PATH=%h/.nix-profile/bin
PassEnvironment=GITHUB_TOKEN FUSION_DASHBOARD_TOKEN FUSION_DAEMON_TOKEN VAST_API_KEY VLLM_API_KEY
WorkingDirectory=%h

RuntimeDirectory=fusion
RuntimeDirectoryMode=0700

ExecStartPre=-%h/.nix-profile/bin/tmux -S %t/fusion/tmux.sock kill-session -t fusion
ExecStart=%h/.nix-profile/bin/fusion-dashboard-service
ExecStop=-%h/.nix-profile/bin/tmux -S %t/fusion/tmux.sock kill-session -t fusion

Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes

ReadOnlyPaths=%h/.npm-packages %h/.nix-profile %h/Sources/Fusion
ReadWritePaths=%h/.fusion %h/.pi %t/fusion %h/.cache %h/.config/qmd

[Install]
WantedBy=default.target
