# .config/systemd/user/fusion.service
[Unit]
Description=Fusion dashboard
Documentation=https://github.com/Runfusion/Fusion
After=network-online.target
Wants=network-online.target

[Service]
PassEnvironment=ANTHROPIC_API_KEY OPENAI_API_KEY OPENROUTER_API_KEY GITHUB_TOKEN FUSION_DASHBOARD_TOKEN FUSION_DAEMON_TOKEN HUGGINGFACE_API_KEY HUGGING_FACE_HUB_TOKEN HF_TOKEN GEMINI_API_KEY GOOGLE_API_KEY
WorkingDirectory=%h/Sources/Fusion
ExecStartPre=-/usr/bin/tmux kill-session -t fusion
ExecStartPre=/usr/bin/mkdir -p %t/fusion
ExecStart=/bin/sh -lc 'rm -f "$XDG_RUNTIME_DIR/fusion/tmux.out"; mkfifo "$XDG_RUNTIME_DIR/fusion/tmux.out"; /usr/bin/tmux new-session -d -s fusion -c "$HOME/Sources/Fusion" "$HOME/.npm-packages/bin/fusion dashboard --host 0.0.0.0 --port 4040"; /usr/bin/tmux pipe-pane -o -t fusion:0.0 "cat > $XDG_RUNTIME_DIR/fusion/tmux.out"; cat "$XDG_RUNTIME_DIR/fusion/tmux.out"'
ExecStop=/usr/bin/tmux kill-session -t fusion
Restart=on-failure
RestartSec=5s
NoNewPrivileges=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ReadWritePaths=%h/Sources/Fusion %h/.fusion %t/fusion

[Install]
WantedBy=default.target
