# .config/systemd/user/home-cleanup.timer
[Unit]
Description=Run home-cleanup weekly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
