# .config/systemd/user/home-cleanup.timer
[Unit]
Description=Run home-cleanup

[Timer]
Persistent=true
OnCalendar=monthly
RandomizedDelaySec=3day

[Install]
WantedBy=timers.target
