# .config/systemd/user/kolu.service
[Unit]
Description=Kolu
Documentation=https://github.com/juspay/kolu
After=network-online.target
Wants=network-online.target

[Service]
Environment=KOLU_OPENCODE_DB=%h/.local/share/opencode/opencode-stable.db
ExecStart=%h/.nix-profile/bin/kolu --host 0.0.0.0
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
