# .config/systemd/user/nix-cleanup.timer
[Unit]
Description=Nix cleanup timer

[Timer]
Persistent=true
OnCalendar=monthly

[Install]
WantedBy=timers.target
