# .config/systemd/user/nix-cleanup.timer
[Unit]
Description=Nix cleanup timer

[Timer]
Persistent=true
OnCalendar=weekly
RandomizedDelaySec=3day

[Install]
WantedBy=timers.target
