# .config/systemd/user/paccache.service
[Unit]
Description=Pacman cache cleanup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pacman-clean
Nice=19

