# .config/systemd/user/paccache.timer
[Unit]
Description=Clean-up old pacman pkg

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3day

[Install]
WantedBy=timers.target
