# .config/systemd/user/pacreport.timer
[Unit]
Description=Pacreport timer

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3day

[Install]
WantedBy=timers.target
