# .config/systemd/user/pacreport.service
[Unit]
Description=Pacreport

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pacman-report
Nice=19

