# .config/systemd/user/vllm-qwen3.6-27B.target
[Unit]
Description=vLLM Qwen3.6 27B MTP + Fusion
Requires=vllm@qwen3.6-27B.service
Wants=fusion.service
Conflicts=vllm-qwen3.6-35B-a3b.target vllm@qwen3.6-35B-a3b.service
After=network-online.target

[Install]
WantedBy=default.target
