# .config/systemd/user/vllm.service
[Unit]
Description=vLLM server via pure Nix runtime derivation
After=default.target network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10
TimeoutStartSec=0
TimeoutStopSec=30
PassEnvironment=HF_TOKEN HF_HUB_TOKEN HUGGING_FACE_HUB_TOKEN VLLM_API_KEY
Environment=MODEL=Intel/Qwen3.6-35B-A3B-int4-mixed-AutoRound
Environment=MODEL_REVISION=65f69c73f17488236c85c85211f6ba28d7106157
Environment=SERVED_MODEL_NAME=qwen3.6-35b-a3b
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=MAX_MODEL_LEN=32768
Environment=MAX_NUM_SEQS=2
Environment=MAX_NUM_BATCHED_TOKENS=3072
Environment=MAX_NUM_PARTIAL_PREFILLS=1
Environment=GPU_MEMORY_UTILIZATION=0.95
Environment=CPU_OFFLOAD_GB=12
Environment=KV_OFFLOADING_SIZE_GB=0
Environment=RENDERER_NUM_WORKERS=4
Environment=MM_PROCESSOR_CACHE_GB=0
Environment=STREAM_INTERVAL=4
Environment=MAX_CUDAGRAPH_CAPTURE_SIZE=4
Environment=SPECULATIVE_CONFIG_B64=eyJtZXRob2QiOiJtdHAiLCJudW1fc3BlY3VsYXRpdmVfdG9rZW5zIjoyfQ==
Environment=VLLM_EXTRA_ARGS=
Environment=OMP_NUM_THREADS=8
Environment=MKL_NUM_THREADS=8
Environment=OPENBLAS_NUM_THREADS=8
Environment=NUMEXPR_NUM_THREADS=8
Environment=RAYON_NUM_THREADS=8
Environment=MALLOC_ARENA_MAX=2
Environment=VLLM_MARLIN_USE_ATOMIC_ADD=1
Environment=HF_TOKEN_FILE=
Environment=HF_HUB_ETAG_TIMEOUT=60
Environment=HF_HUB_DOWNLOAD_TIMEOUT=600
Environment=HF_HUB_DISABLE_IMPLICIT_TOKEN=0
Environment=SHOW_HIDDEN_METRICS_FOR_VERSION=
Environment=ENABLE_LOGGING_ITERATION_DETAILS=0
Environment=DOWNLOAD_STATUS_AT_START=1
Environment=BENCHMARK_PROMPTS="1024 4096 12000"
Environment=BENCHMARK_MAX_TOKENS="96 128"
Environment=VLLM_NIX_CONFIG=%h/.nix/vllm-runtime.nix
Environment=VLLM_RUNTIME_DIR=%C/vllm-nix-qwen36

UnsetEnvironment=LD_LIBRARY_PATH LD_PRELOAD

CacheDirectory=vllm-nix-qwen36
CacheDirectoryMode=0700

ExecStart=/bin/sh -lc 'set -eu; nix_bin=""; for candidate in "$HOME/.nix-profile/bin/nix" /nix/var/nix/profiles/default/bin/nix /run/current-system/sw/bin/nix /usr/bin/nix; do if [ -x "$candidate" ]; then nix_bin="$candidate"; break; fi; done; if [ -z "$nix_bin" ]; then echo "nix executable not found" >&2; exit 127; fi; CACHE="$VLLM_RUNTIME_DIR"; export XDG_CACHE_HOME="$CACHE/xdg-cache"; mkdir -p "$CACHE" "$XDG_CACHE_HOME"; VLLM_STORE_PATH="$($nix_bin --extra-experimental-features "nix-command flakes" build --no-link --print-out-paths --file "$VLLM_NIX_CONFIG")"; if [ -r "$HOME/.profile" ]; then set -a; . "$HOME/.profile"; [ ! -r "$HOME/.profile_override" ] || . "$HOME/.profile_override"; set +a; fi; TOKEN_SOURCE="environment-or-profile"; if [ -z "${HF_TOKEN:-}" ] && [ -n "${HF_HUB_TOKEN:-}" ]; then export HF_TOKEN="$HF_HUB_TOKEN"; TOKEN_SOURCE="HF_HUB_TOKEN"; fi; if [ -z "${HF_TOKEN:-}" ] && [ -n "${HUGGING_FACE_HUB_TOKEN:-}" ]; then export HF_TOKEN="$HUGGING_FACE_HUB_TOKEN"; TOKEN_SOURCE="HUGGING_FACE_HUB_TOKEN"; fi; if [ -z "${HF_TOKEN:-}" ]; then for token_file in "${HF_TOKEN_FILE:-}" "$HOME/.cache/huggingface/token" "$HOME/.huggingface/token"; do if [ -n "$token_file" ] && [ -r "$token_file" ]; then HF_TOKEN="$(head -n 1 "$token_file")"; export HF_TOKEN; TOKEN_SOURCE="file:$token_file"; break; fi; done; fi; export HF_HUB_TOKEN="${HF_HUB_TOKEN:-${HF_TOKEN:-}}"; export HUGGING_FACE_HUB_TOKEN="${HUGGING_FACE_HUB_TOKEN:-${HF_TOKEN:-}}"; export HF_TOKEN="${HF_TOKEN:-${HF_HUB_TOKEN:-${HUGGING_FACE_HUB_TOKEN:-}}}"; export HF_HUB_ETAG_TIMEOUT="${HF_HUB_ETAG_TIMEOUT:-60}"; export HF_HUB_DOWNLOAD_TIMEOUT="${HF_HUB_DOWNLOAD_TIMEOUT:-600}"; export HF_HUB_DISABLE_IMPLICIT_TOKEN="${HF_HUB_DISABLE_IMPLICIT_TOKEN:-0}"; if [ -n "${HF_TOKEN:-}" ]; then echo "HF_TOKEN is set for vLLM/Hugging Face access via $TOKEN_SOURCE"; else echo "WARNING: HF_TOKEN is not set after profile and token-file checks"; fi; CUDA_SOURCE_DIR=""; LIBCUDA_SO=""; for d in /run/opengl-driver/lib /run/nvidia-driver/lib /usr/lib/x86_64-linux-gnu/nvidia/current /usr/lib/x86_64-linux-gnu /usr/lib64 /usr/lib /usr/lib/nvidia; do if [ -e "$d/libcuda.so.1" ]; then CUDA_SOURCE_DIR="$d"; LIBCUDA_SO="$d/libcuda.so.1"; break; fi; done; if [ -z "$LIBCUDA_SO" ]; then LIBCUDA_SO="$(find /run /usr/lib /usr/lib64 -maxdepth 6 -name libcuda.so.1 2>/dev/null | head -n 1 || true)"; [ -z "$LIBCUDA_SO" ] || CUDA_SOURCE_DIR="$(dirname "$LIBCUDA_SO")"; fi; if [ -z "$LIBCUDA_SO" ] || [ ! -e "$LIBCUDA_SO" ]; then echo "ERROR: libcuda.so.1 not found in common driver paths"; find /run /usr/lib /usr/lib64 -maxdepth 6 -name libcuda.so.1 2>/dev/null | head -50 || true; exit 1; fi; DRIVER_LINK_DIR="$CACHE/driver-libs"; rm -rf "$DRIVER_LINK_DIR"; mkdir -p "$DRIVER_LINK_DIR"; for f in "$CUDA_SOURCE_DIR"/libcuda.so* "$CUDA_SOURCE_DIR"/libnvidia-*.so*; do [ -e "$f" ] || continue; ln -sfn "$(readlink -f "$f")" "$DRIVER_LINK_DIR/$(basename "$f")"; done; [ -e "$DRIVER_LINK_DIR/libcuda.so.1" ] || ln -sfn "$(readlink -f "$LIBCUDA_SO")" "$DRIVER_LINK_DIR/libcuda.so.1"; [ -e "$DRIVER_LINK_DIR/libcuda.so" ] || ln -sfn libcuda.so.1 "$DRIVER_LINK_DIR/libcuda.so"; export HF_HOME="${HF_HOME:-$HOME/.cache/huggingface}"; export VLLM_USE_V1=1; export VLLM_SLEEP_WHEN_IDLE=1; export CUDA_MODULE_LOADING="${CUDA_MODULE_LOADING:-LAZY}"; export VLLM_MARLIN_USE_ATOMIC_ADD="${VLLM_MARLIN_USE_ATOMIC_ADD:-1}"; export PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True; export OMP_NUM_THREADS="${OMP_NUM_THREADS:-8}"; export MKL_NUM_THREADS="${MKL_NUM_THREADS:-8}"; export OPENBLAS_NUM_THREADS="${OPENBLAS_NUM_THREADS:-8}"; export NUMEXPR_NUM_THREADS="${NUMEXPR_NUM_THREADS:-8}"; export RAYON_NUM_THREADS="${RAYON_NUM_THREADS:-8}"; export TOKENIZERS_PARALLELISM="${TOKENIZERS_PARALLELISM:-true}"; export MALLOC_ARENA_MAX="${MALLOC_ARENA_MAX:-2}"; PY_SITE="$(find "$VLLM_STORE_PATH/lib" -type d -name site-packages -print 2>/dev/null | head -n 1 || true)"; if [ -z "$PY_SITE" ]; then echo "ERROR: Python site-packages directory not found under $VLLM_STORE_PATH/lib"; find "$VLLM_STORE_PATH/lib" -maxdepth 4 -type d 2>/dev/null | head -50 || true; exit 1; fi; TORCH_LIB_DIR="$PY_SITE/torch/lib"; CUDA_PKG_LIB_DIRS=""; if [ -d "$PY_SITE/nvidia" ]; then CUDA_PKG_LIB_DIRS="$(find "$PY_SITE/nvidia" -type d -name lib -print 2>/dev/null | paste -sd: - || true)"; fi; RUNTIME_LD_PATH="$DRIVER_LINK_DIR"; [ ! -d "$TORCH_LIB_DIR" ] || RUNTIME_LD_PATH="$TORCH_LIB_DIR:$RUNTIME_LD_PATH"; [ -z "$CUDA_PKG_LIB_DIRS" ] || RUNTIME_LD_PATH="$CUDA_PKG_LIB_DIRS:$RUNTIME_LD_PATH"; export LD_LIBRARY_PATH="$RUNTIME_LD_PATH:${LD_LIBRARY_PATH:-}"; RUNTIME_ARGS=(--served-model-name "$SERVED_MODEL_NAME" --host "$HOST" --port "$PORT" --trust-remote-code --dtype half --max-model-len "$MAX_MODEL_LEN" --max-num-seqs "$MAX_NUM_SEQS" --gpu-memory-utilization "$GPU_MEMORY_UTILIZATION" --cpu-offload-gb "$CPU_OFFLOAD_GB" --max-num-batched-tokens "$MAX_NUM_BATCHED_TOKENS" --max-num-partial-prefills "$MAX_NUM_PARTIAL_PREFILLS" --renderer-num-workers "$RENDERER_NUM_WORKERS" --mm-processor-cache-gb "$MM_PROCESSOR_CACHE_GB" --stream-interval "$STREAM_INTERVAL" --max-cudagraph-capture-size "$MAX_CUDAGRAPH_CAPTURE_SIZE" --enable-chunked-prefill --enable-prefix-caching --async-scheduling --language-model-only --skip-mm-profiling --reasoning-parser qwen3 --enable-auto-tool-choice --tool-call-parser qwen3_coder --disable-uvicorn-access-log); [ -z "${MODEL_REVISION:-}" ] || RUNTIME_ARGS+=(--revision "$MODEL_REVISION" --tokenizer-revision "$MODEL_REVISION" --code-revision "$MODEL_REVISION"); [ -z "${VLLM_API_KEY:-}" ] || RUNTIME_ARGS+=(--api-key "$VLLM_API_KEY"); if [ -n "${KV_OFFLOADING_SIZE_GB:-}" ] && [ "${KV_OFFLOADING_SIZE_GB:-0}" != "0" ]; then RUNTIME_ARGS+=(--kv-offloading-size "$KV_OFFLOADING_SIZE_GB"); fi; if [ -n "${SPECULATIVE_CONFIG_B64:-}" ]; then SPEC_JSON="$(printf %s "$SPECULATIVE_CONFIG_B64" | base64 -d)"; RUNTIME_ARGS+=(--speculative-config "$SPEC_JSON"); fi; [ "${ENABLE_LOGGING_ITERATION_DETAILS:-0}" != "1" ] || RUNTIME_ARGS+=(--enable-logging-iteration-details); [ -z "${SHOW_HIDDEN_METRICS_FOR_VERSION:-}" ] || RUNTIME_ARGS+=(--show-hidden-metrics-for-version "$SHOW_HIDDEN_METRICS_FOR_VERSION"); if [ -n "${VLLM_EXTRA_ARGS:-}" ]; then read -r -a EXTRA_ARGS <<< "$VLLM_EXTRA_ARGS"; RUNTIME_ARGS+=("${EXTRA_ARGS[@]}"); fi; printf "%q " "$VLLM_STORE_PATH/bin/vllm" serve "$MODEL" "${RUNTIME_ARGS[@]}" > "$CACHE/runtime_command.sh"; echo >> "$CACHE/runtime_command.sh"; chmod +x "$CACHE/runtime_command.sh"; echo "Runtime command: $CACHE/runtime_command.sh"; echo "Using pure Nix vLLM runtime: $VLLM_STORE_PATH"; echo "Using CUDA driver source: $LIBCUDA_SO"; echo "Using isolated driver lib dir: $DRIVER_LINK_DIR"; echo "Using Python site-packages: $PY_SITE"; echo "Using CUDA package lib dirs: ${CUDA_PKG_LIB_DIRS:-none}"; echo "Using torch lib dir: ${TORCH_LIB_DIR:-none}"; echo "Serving model: $MODEL as $SERVED_MODEL_NAME"; exec "$VLLM_STORE_PATH/bin/vllm" serve "$MODEL" "${RUNTIME_ARGS[@]}"'

ReadOnlyPaths=%h/.nix/vllm-runtime.nix
ReadWritePaths=%C/vllm-nix-qwen36 %h/.cache/huggingface

[Install]
WantedBy=default.target
