# .config/systemd/user/vllm@.service
[Unit]
Description=vLLM model server (%i)
After=default.target network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10
TimeoutStartSec=0
TimeoutStopSec=30
PassEnvironment=HF_TOKEN HF_HUB_TOKEN HUGGING_FACE_HUB_TOKEN VLLM_API_KEY
EnvironmentFile=%h/.config/vllm/%i.env
Environment=VLLM_INSTANCE=%i
Environment=VLLM_NIX_CONFIG=%h/.config/nix/vllm-runtime.nix
Environment=VLLM_RUNTIME_DIR=%C/vllm-%i
Environment=HF_HUB_ETAG_TIMEOUT=60
Environment=HF_HUB_DOWNLOAD_TIMEOUT=600
Environment=HF_HUB_DISABLE_IMPLICIT_TOKEN=0
Environment=FLASHINFER_CACHE_DIR=%C/vllm-%i/flashinfer
Environment=FLASHINFER_WORKSPACE_BASE=%C/vllm-%i/flashinfer
Environment=TRITON_CACHE_DIR=%C/vllm-%i/triton-cache
Environment=TRITON_HOME=%C/vllm-%i/triton-home
Environment=VLLM_NO_USAGE_STATS=1
Environment=DO_NOT_TRACK=1
Environment=OMP_NUM_THREADS=8
Environment=MKL_NUM_THREADS=8
Environment=OPENBLAS_NUM_THREADS=8
Environment=NUMEXPR_NUM_THREADS=8
Environment=RAYON_NUM_THREADS=8
Environment=MALLOC_ARENA_MAX=2
Environment=VLLM_MARLIN_USE_ATOMIC_ADD=1
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True
Environment=VLLM_READY_TIMEOUT_SECONDS=1800

UnsetEnvironment=LD_LIBRARY_PATH LD_PRELOAD

CacheDirectory=vllm-%i
CacheDirectoryMode=0700

ExecStartPre=%h/.local/bin/vllm-patch-model-defaults
ExecStart=%h/.local/bin/vllm-serve-pure
ExecStartPost=%h/.local/bin/vllm-wait-ready
ExecStartPost=/usr/bin/systemctl --user restart fusion.service

NoNewPrivileges=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes

ReadOnlyPaths=%h/.config/nix/vllm-runtime.nix %h/.config/vllm/%i.env %h/.local/bin/vllm-patch-model-defaults %h/.local/bin/vllm-serve-pure %h/.local/bin/vllm-wait-ready
ReadWritePaths=%C/vllm-%i %h/.cache/huggingface %h/.config/crush %h/.pi/agent %h/.fusion

[Install]
WantedBy=default.target
